# conf/cron/services_docker.cron
* * * * * . /etc/profile.d/02_env.sh; python3 /dmarket/manage.py runcrons >> /tmp/cron.log 2>&1
