# conf/cron/services.cron
* * * * * . /etc/profile.d/02_env.sh; python3 $DMARKET_ROOT/manage.py runcrons > /tmp/cron2.log 2>&1
